# scripts/JF01-buffer.service
[Unit]
Description=All UDP-buffer instances of JF01

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /home/dbe/service_scripts/JF01-buffer-worker.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
