# scripts/JF11-buffer.service
[Unit]
Description=All UDP-buffer instances of JF11

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /home/dbe/service_scripts/JF11-buffer-worker.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
