# scripts/JF11-stream.service
[Unit]
Description=stream service (to streamvis and live analysis) of JF11(TXS Flex)

[Service]
PermissionsStartOnly=true
Type=idle
User=root
ExecStart=/usr/bin/sh /home/dbe/service_scripts/JF11-stream.sh
TimeoutStartSec=10
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target

