# scripts/JF13-buffer.service
[Unit]
Description=All UDP-buffer instances of JF13

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /home/writer/git/sf_daq_buffer/scripts/JF13-buffer-worker.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
